# s6-rc.d/freshclam/run
#!/command/execlineb -P

exec freshclam --daemon
